# system/etc/udev/rules.d/90-mouse.rules
SUBSYSTEM=="input", KERNEL=="mouse[0-9]*", ACTION=="add|remove", RUN+="/usr/local/bin/pick-peripheral-udev mouse"
